# systemd/openclaw-avatar-status.service
[Unit]
Description=OpenClaw avatar HTTP status feed
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/pi-avatar/avatar.env
WorkingDirectory=/opt/pi-avatar
ExecStart=/opt/pi-avatar/.venv/bin/python /opt/pi-avatar/status_agent.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
