# systemd/pi-avatar-renderer.service
[Unit]
Description=Pi Avatar renderer
After=local-fs.target

[Service]
Type=simple
EnvironmentFile=-/etc/pi-avatar/avatar.env
WorkingDirectory=/opt/pi-avatar
ExecStart=/opt/pi-avatar/.venv/bin/python /opt/pi-avatar/renderer.py
Restart=always
RestartSec=2
StandardInput=tty
TTYPath=/dev/tty1

[Install]
WantedBy=multi-user.target
