# systemd/pi-avatar-monitor.service
[Unit]
Description=Pi Avatar remote OpenClaw monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/pi-avatar/avatar.env
WorkingDirectory=/opt/pi-avatar
ExecStart=/usr/bin/python3 /opt/pi-avatar/monitor.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
